# scripts/crontab
# CrowdSec → feed + MISP (existant)
0 1,13 * * * python /app/feed.py

# MISP event → feed GitHub (nouveau) — décalé de 30 min
30 1,13 * * * python /app/misp_export.py
